# other/systemd/dbhub-webui.socket
[Unit]
Description=Socket that allows the DBHub.io WebUI to bind to port 443

[Socket]
ListenStream=443
NoDelay=true
